# dist/systemd/wgconfd-state@.service
# SPDX-License-Identifier: CC0-1.0
#
# Copyright 2019 Hristo Venev

[Unit]
Description=WireGuard runtime state for %i

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=true
RuntimeDirectory=wgconfd/%i
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=no
