# dist/systemd/wgconfd@.service
[Unit]
Description=WireGuard configuration daemon on %i
Wants=network.target
After=network-pre.target
Before=network.target

[Service]
Type=simple
CapabilityBoundingSet=CAP_NET_ADMIN
Restart=on-failure
RestartSec=0
RuntimeDirectory=wgconfd/%i
RuntimeDirectoryPreserve=yes
CacheDirectory=wgconfd/%i
ExecStart=/usr/bin/env wgconfd %i /etc/wireguard/%i.toml
StandardError=journal
SyslogLevelPrefix=true

[Install]
WantedBy=multi-user.target
